# celery.service
# /home/ubuntu/stable-diffusion/celery.service => /etc/systemd/system/celery.service
[Unit]
Description=Celery daemon
After=redis.service

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/stable-diffusion
ExecStart=/home/ubuntu/anaconda3/envs/ldm/bin/celery -A scripts.sdserver.celery worker --loglevel=INFO
Restart=always

[Install]
WantedBy=multi-user.target